# rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
with-contenv
s6-setuidgid rtorrent
nginx -g "daemon off;"